# crontab
25 14 * * * bash  ~/SAE51_2/script/save.sh
